# powershell/.config/powershell/profile.d/00-modules.ps1
#TODO: Add verifications to warn that this modules not installed
Import-Module -Name Terminal-Icons
Import-Module -Name PSReadLine
