# powershell/.config/powershell/profile.d/10-aliases.ps1
Set-Alias c clear
Set-Alias ff fastfetch
