# powershell/.config/powershell/profile.d/40-utils.ps1
#TODO: Add verifications to warn that this modules not installed
Invoke-Expression (& { (zoxide init powershell | Out-String) }) # zoxide init
